# python/flask_app.conf
[supervisord]
nodaemon=true

[program:flask_app]
command=python3 /home/geremytanyentsen/flask_api.py
directory=/home/geremytanyentsen
autostart=true
autorestart=true